# crontab
SHELL=/bin/bash
PATH=/app/venv/bin:/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 0 * * * /app/venv/bin/python /app/main.py >> /var/log/cron.log 2>&1